# meta-inventec/meta-c600g5/recipes-phosphor/gpio-monitor/c600g5-gpio-monitor-register/cpu-throttle-deassert.service
[Unit]
Description=c600g5 CPU Throttle event deassert Service

[Service]
Type=oneshot
ExecStart=/usr/bin/env identify_cpu_throttle_event.sh throttle skip deassert
