# meta-inventec/meta-c600g5/recipes-phosphor/gpio-monitor/c600g5-gpio-monitor-register/power-button-released.service
[Unit]
Description=c600g5 Power Button Released Service

[Service]
Type=oneshot
ExecStart=/usr/bin/env power_button_event.sh released
